# config/systemd/user/change_wallpaper.timer
[Unit]
Description=Timer to change wallpper

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
